# DEBIAN/postinst
#!/bin/bash
set -e

## Backup resolv.conf
if [ -f /etc/resolv.conf ]; then
    cp /etc/resolv.conf /etc/resolv.conf.backup || true
fi

## Fix permissions and ownership of resolv.conf
chown root:root /etc/resolv.conf
chmod 644 /etc/resolv.conf

## Copy privacy-conf dnscrypt-proxy settings
if [ -f /etc/dnscrypt-proxy/dnscrypt-sirros.toml ]; then
    cp /etc/dnscrypt-proxy/dnscrypt-sirros.toml /etc/dnscrypt-proxy/dnscrypt-proxy.toml
else
    echo "Warning: dnscrypt-sirros.toml not found, skipping."
fi

## Comment out the User line in the systemd unit
UNIT_FILE="/etc/systemd/system/dnscrypt-proxy.service"
if [ -f "$UNIT_FILE" ]; then
    sed -i 's/^[[:space:]]*User=.*/#&/' "$UNIT_FILE"
fi

# Stop systemd-resolved and mask its sockets to free port 53
if systemctl list-unit-files | grep -q systemd-resolved; then
    systemctl stop systemd-resolved || true
    systemctl stop systemd-resolved-varlink.socket || true
    systemctl stop systemd-resolved-monitor.socket || true

    systemctl disable systemd-resolved || true
    systemctl mask systemd-resolved || true
    systemctl mask systemd-resolved-varlink.socket || true
    systemctl mask systemd-resolved-monitor.socket || true

## Reload systemd to pick up changes
systemctl daemon-reexec

## Restart NetworkManager to apply configs
if systemctl is-active --quiet NetworkManager; then
    systemctl restart NetworkManager
fi

## Restart dnscrypt-proxy to apply changes
systemctl restart dnscrypt-proxy

## Optional: make resolv.conf immutable
# chattr +i /etc/resolv.conf || true

echo "privacy-conf installation complete: DNS is forced to localhost (127.0.0.1) and dnscrypt-proxy runs without User restriction."
